# services/ffmpegcmd/finish
#!/bin/bash
echo 'ffmpeg command exited... restarting'
