# services/ffmpegcmd/run
#!/usr/bin/with-contenv /bin/bash
echo 'Starting ffmpeg'
exec python /root/runffmpeg.py

