# services/ffmpegcmd/run
#!/usr/bin/with-contenv /bin/bash
echo 'Running ffmpeg command'
IFS=', ' read -r -a array <<< "$string"
echo $CMD
exec ffmpeg "${CMD[@]}"

