# services/cron/finish
#!/bin/bash
echo cron exited... restarting

